# config/systemd/user/AlesisQ88Transport.service
[Unit]
Description=Alesis Q88 Transport to Mackie

[Service]
ExecStart=/usr/bin/AlesisQ88_transport

[Install]
WantedBy=default.target
