# config/systemd/user/dropbox.service.d/override.conf
[Service]
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
